# scripts/systemd/easyproxy.service
[Unit]
Description=EasyProxy Multi-Protocol Proxy Server
Documentation=https://github.com/Slothtron/easyproxy-py
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# 运行用户和组 (建议创建专用用户)
User=easyproxy
Group=easyproxy

# 工作目录
WorkingDirectory=/opt/easyproxy

# 启动命令 (根据实际安装路径调整)
ExecStart=/usr/bin/easyproxy start -c /etc/easyproxy/config.yaml

# 优雅重载配置 (发送 SIGHUP)
ExecReload=/bin/kill -HUP $MAINPID

# 优雅停止 (发送 SIGTERM)
ExecStop=/bin/kill -TERM $MAINPID

# 自动重启策略
Restart=always
RestartSec=10

# 重启限制 (200秒内最多重启5次,防止重启风暴)
StartLimitInterval=200
StartLimitBurst=5

# 日志输出到 systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=easyproxy

# 环境变量
Environment="PYTHONUNBUFFERED=1"

# ===== 性能优化 =====

# 文件描述符限制 (支持大量并发连接)
LimitNOFILE=1048576

# 进程数限制
LimitNPROC=unlimited

# 核心转储大小限制 (用于调试崩溃)
LimitCORE=infinity

# ===== 安全加固 =====

# 禁止提升权限
NoNewPrivileges=true

# 私有临时目录
PrivateTmp=true

# 保护系统目录 (只读)
ProtectSystem=strict

# 保护 /home 目录
ProtectHome=true

# 只读路径
ReadOnlyPaths=/

# 可写路径 (日志目录)
ReadWritePaths=/var/log/easyproxy

# 隔离设备访问
PrivateDevices=true

# 限制系统调用 (只允许网络服务相关的系统调用)
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# 能力限制 (只保留绑定端口的权限)
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# 禁止访问内核模块
ProtectKernelModules=true

# 禁止访问内核日志
ProtectKernelLogs=true

# 禁止修改内核参数
ProtectKernelTunables=true

# 保护控制组
ProtectControlGroups=true

# 限制地址族 (只允许 IPv4 和 IPv6)
RestrictAddressFamilies=AF_INET AF_INET6

# 锁定个性化设置
LockPersonality=true

# 限制实时调度
RestrictRealtime=true

# 禁止 SUID/SGID
RestrictSUIDSGID=true

# 移除 IPC 命名空间
PrivateIPC=true

# ===== 资源限制 (可选,根据实际需求调整) =====

# CPU 配额 (限制为 2 个核心)
# CPUQuota=200%

# 内存限制 (限制为 2GB)
# MemoryLimit=2G

# 任务数限制
# TasksMax=4096

[Install]
WantedBy=multi-user.target
